# berks-cookbooks/redisio/templates/default/redis@.service.erb
[Unit]
Description=Redis (%i) persistent key-value database
After=network.target

[Service]
ExecStart=<%= @bin_path %>/redis-server /etc/redis/%i.conf --daemonize no
User=<%= @user %>
Group=<%= @group %>
LimitNOFILE=<%= @limit_nofile %>

[Install]
WantedBy=multi-user.target
